# salt_highstate.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

*/15 * * * * root salt-run state.orch orch.highstate 
